# deploy/systemd/veloxhash.service
[Unit]
Description=VeloxHash miner
Documentation=file:/etc/veloxhash/README.systemd
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=exec
User=veloxhash
Group=veloxhash
EnvironmentFile=/etc/veloxhash/veloxhash.env
UMask=0027
StandardInput=null
StandardOutput=journal
StandardError=journal
RuntimeDirectory=veloxhash
RuntimeDirectoryMode=0750
StateDirectory=veloxhash
StateDirectoryMode=0750
LogsDirectory=veloxhash
LogsDirectoryMode=0750
WorkingDirectory=/var/lib/veloxhash
ExecStart=/usr/local/lib/veloxhash/veloxhash-run
Restart=always
RestartSec=10
KillSignal=SIGTERM
TimeoutStopSec=30
Nice=10
CPUWeight=20
IOWeight=20
IOSchedulingClass=idle
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/etc/veloxhash /var/lib/veloxhash /var/log/veloxhash
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
